# roles/transplant/templates/upstart.conf.j2
start on runlevel [2345]
stop on runlevel [06]
respawn

setuid {{ transplant.user }}
setgid {{ transplant.group }}

exec {{ transplant.dir }}/env/bin/uwsgi --ini {{ transplant.dir }}/uwsgi.ini
